# deploy/systemd/youtube-digest.service
[Unit]
Description=Generate YouTube Digest EPUB
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/youtube-digest
EnvironmentFile=/opt/youtube-digest/.env
ExecStart=/opt/youtube-digest/.venv/bin/youtube-digest --config /opt/youtube-digest/config.json generate --json
StandardOutput=append:/opt/youtube-digest/logs/youtube-digest.log
StandardError=append:/opt/youtube-digest/logs/youtube-digest-error.log
